# assets/agent-stack/openclaw.container
[Unit]
Description=openclaw agent (upstream NemoClaw OSS)
Documentation=https://github.com/openclaw/openclaw
Wants=network-online.target nclawzero-load-agent-images.service
After=network-online.target nclawzero-load-agent-images.service

[Container]
Image=ghcr.io/openclaw/openclaw@sha256:06b4f3dfa3c88d49c92e99d635dc62053d4afd045d6220e811dff6190040f3de
# Image pinned to immutable digest. To roll forward, resolve a newer
# sha256 via `crane digest ghcr.io/openclaw/openclaw:main` and update
# here. AutoUpdate=registry intentionally OMITTED to prevent unreviewed
# upstream changes from running with API keys.
ContainerName=openclaw
EnvironmentFile=/etc/nclawzero/agent-env
Volume=openclaw-data:/openclaw-data:Z
Volume=/var/lib/nclawzero/openclaw-home:/home/node/.openclaw:Z
PublishPort=18789:18789

[Service]
TimeoutStartSec=600
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target default.target
