# assets/agent-stack/zeroclaw.container
[Unit]
Description=zeroclaw agent gateway (nclawzero stack)
Documentation=https://gitlab.com/nclawzero/zeroclaw
Wants=network-online.target nclawzero-load-agent-images.service
After=network-online.target nclawzero-load-agent-images.service

[Container]
Image=ghcr.io/perlowja/nclawzero-demo@sha256:0d5306ff1fb35602dfb283bf35390bb88b91e5164a511d1e2af5800d48f627c0
# Image pinned to immutable digest. To roll forward to a newer build,
# resolve the new sha256 via `crane digest ghcr.io/perlowja/nclawzero-
# demo:latest` (or the curl manifest pattern in nclawzero-agent-quadlets
# bb), update this line, and ship a new recipe rev. AutoUpdate=registry
# intentionally OMITTED to prevent unreviewed container code from
# running with API keys.
ContainerName=zeroclaw
Entrypoint=/usr/local/bin/zeroclaw
Exec=daemon --host 0.0.0.0 --port 42617
EnvironmentFile=/etc/nclawzero/agent-env
Volume=zeroclaw-data:/zeroclaw-data:Z
PublishPort=42617:42617

[Service]
TimeoutStartSec=900
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target default.target
